# share/asmd_smf.xml.in
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='asmd'>
  <service name='system/asmd' type='service' version='0'>
{instances}
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='python'>Advanced SmartOS Management Daemon</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
